# infra/supervisor/supervisord.conf
; IoTAPS background worker supervision (Req 30.1, 30.3)
; Runs all 8 worker processes. Each program has autorestart=true so a worker
; that terminates unexpectedly is restarted immediately. Combined with the
; container's `restart: always`, this gives two layers of self-healing.

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:mqtt-listener]
command=python -m app.workers.mqtt_listener
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:batch-writer]
command=python -m app.workers.batch_writer
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:downsampler]
command=python -m app.workers.downsampler
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:alert-checker]
command=python -m app.workers.alert_checker
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:webhook-dispatcher]
command=python -m app.workers.webhook_dispatcher
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:notification-sender]
command=python -m app.workers.notification_sender
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:session-cleanup]
command=python -m app.workers.session_cleanup
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:data-retention]
command=python -m app.workers.data_retention
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
